# fs/etc/services.d/nordvpn/finish
#!/bin/bash

rm -rf /run/nordvpn
ip link delete nordlynx 2>/dev/null